# mern/supervisor/supervisord.conf
[supervisord]
nodaemon=false

[program:mongodb]
command=/usr/bin/mongod "mongod --dbpath /usr/data/db"

[program:process1]
command=/bin/bash "cd /usr/src/app && npm start"
